# cloudcmd/rootfs/etc/s6-overlay/s6-rc.d/init-cloudcmd/run
#!/usr/bin/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================

if ! bashio::fs.file_exists '/config/settings.json'; then
  echo "{}" > /config/settings.json
fi
